# cron/cron.d/website-update
*/5 * * * * www $HOME/website/cron/scripts/website-update.sh
